# f.conf
[program:f]
command=python3 /home/curator2/f/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/f/test.err.log
stdout_logfile=/home/curator2/f/test.out.log